# crontab
# m h  dom mon dow   command
*/5 * * * * python ~/mengivakt/thingspeak.py
*/1 * * * * ~/mengivakt/respawnIfNeeded.sh
*/1 * * * * ~/mengivakt/pi/cam/isNighttime.py

